# deploy/findata-br.service
[Unit]
Description=findata-br — open-source Brazilian financial data MCP server
Documentation=https://github.com/robertoecf/findata-br
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=findata
Group=findata
WorkingDirectory=/opt/findata-br
Environment="PATH=/opt/findata-br/.venv/bin:/usr/local/bin:/usr/bin:/bin"
# Rate-limit tuning for public deploys. Override in a drop-in override.conf
# if you want tighter/looser buckets per deployment.
Environment="FINDATA_RATE_LIMIT_DEFAULT=60/minute;1000/day"
Environment="FINDATA_RATE_LIMIT_ENABLED=true"
ExecStart=/opt/findata-br/.venv/bin/findata serve --host 0.0.0.0 --port 8000 --no-banner
Restart=on-failure
RestartSec=5
# Security hardening — read-only service talking to public APIs only.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/findata-br
ProtectKernelTunables=true
ProtectControlGroups=true
MemoryMax=512M
TasksMax=256

[Install]
WantedBy=multi-user.target
